# src/main/crontab
#START CRON
15,45 * * * * python /app/main.py
#END CRON